# systemd/magic-mouse-gnome.service
[Unit]
Description=Magic Mouse 2 Gesture Driver (GNOME/KDE Wayland)
Documentation=https://github.com/longstoryscott/magic-mouse-gnome
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/magic-mouse-gnome/magic_mouse_gestures.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
